# systemd/tm-monitor.service
[Unit]
Description=Tunnel Manager - health monitor & auto-recovery
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tunnelctl monitor-run
Restart=always
RestartSec=5
Nice=5

[Install]
WantedBy=multi-user.target
