# systemd/tm-report.service
[Unit]
Description=Tunnel Manager - scheduled daily report
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tunnelctl report daily --send
